# bin/wetty.service
# systemd unit file
#
# place in /etc/systemd/system
# systemctl enable wetty.service
# systemctl start wetty.service

[Unit]
Description=Wetty Web Terminal
After=network.target

[Service]
User=root
Group=root

WorkingDirectory=/home/akadmin/pack/wetty # CHANGE ME
ExecStart=/usr/bin/node app.js -p 9123 --host 127.0.0.1

[Install]
WantedBy=multi-user.target
