# bin/wetty.conf
# Upstart script
# /etc/init/wetty.conf

description "Web TTY"
author      "Wetty"

start on started network
stop on stopping network

respawn
respawn limit 20 5

exec wetty -p 3000
